# roles/tc/templates/tc.service.j2
[Unit]
Description=Configure TC settings
After=network.target

[Service]
ExecStart=/bin/bash {{ config_dir }}/tc.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
